# .config/autostart/touchpad-tweaks.desktop
[Desktop Entry]
Type=Application
Exec=/home/mkelly/bin/touchpad-tweaks.sh
